# bifrost.service
[Unit]
Description=Bifröst Gate - VPN Monitoring Agent
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/etc/bifrost
ExecStart=/usr/bin/bifrost-gate
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
