# debian/postinst
#!/bin/sh
set -e

echo "=== Bifröst-Gate: Script de Post-instalación ==="

# 1. Crear directorios mínimos de configuración
echo "Verificando directorios..."
mkdir -p /etc/bifrost/certs
chmod 755 /etc/bifrost

# 4. Verificar logrotate (no instalar desde postinst)
echo "Verificando logrotate..."
if ! command -v logrotate >/dev/null 2>&1; then
    echo "⚠ Advertencia: 'logrotate' no está instalado."
    echo "  Recomendado: instálalo con tu gestor de paquetes (p.ej. 'apt-get install logrotate')."
    echo "  La rotación de logs no estará activa hasta que lo instales."
fi

# 5. Aplicar tmpfiles.d (logs, /var/lib/bifrost) si systemd-tmpfiles está disponible
if command -v systemd-tmpfiles >/dev/null 2>&1; then
    systemd-tmpfiles --create /usr/lib/tmpfiles.d/bifrost-gate.conf 2>/dev/null || true
fi

# 7. Recargar systemd
echo "Recargando configuración de systemd..."
systemctl daemon-reload

# 8. Habilitar el servicio (pero no iniciarlo)
systemctl enable bifrost-gate.service 2>/dev/null || true

# 9. Verificar permisos de config.toml
if [ -f /etc/bifrost/config.toml ]; then
    chmod 600 /etc/bifrost/config.toml
    echo "Permisos de config.toml ajustados (600)."
fi

echo "✓ Bifröst-Gate instalado correctamente."
echo "  Archivos de logs: /var/log/bifrost/*.log"
echo "  Directorio de datos: /var/lib/bifrost"
echo "  Directorio de configuración: /etc/bifrost"
echo ""
echo "Para iniciar el servicio, ejecuta:"
echo "  systemctl start bifrost-gate"
echo ""
echo "Para ver los logs:"
echo "  journalctl -u bifrost-gate -f"
echo "  O: tail -f /var/log/bifrost/*.log"
echo ""

